# config/autostart/mouse.desktop
[Desktop Entry]
Type=Application
Name=mouse
Comment=slower mouse script
Exec=/home/tan/.config/mouse.sh
Hidden=true

